# check_link.service
[Unit]
Description=Check connection to Server

[Service]
Type=simple
ExecStart=/opt/killmywatt/check_link.sh
